# test.timer
[Unit]
Description=Timer for monitor process test

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m
AccuracySec=1s

[Install]
WantedBy=timers.target
